# redis/multiport/redis_18001.conf
[program:redis_18001]
#command=/etc/init.d/redis-server start 18001
command=/usr/bin/redis-server /etc/redis/18001.conf
process_name=redis_18001
autorestart=true
redirect_stderr=true
startsecs=0
stderr_logfile=/data/logs/redis_18001.err.log
stdout_logfile=/data/logs/redis_18001.out.log
directory=/etc/redis
user=root
